# resources/systemd/backup.service
[Unit]
Description=Backup of cassandra database
After=multi-user.target
Wants=multi-user.target
[Service]
Type=oneshot
Environment="PATH=/opt/mesosphere/bin:/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin:/opt/bin"
ExecStart=/home/core/cassandra-backup/backup-witan.sh
