# resources/systemd/cadvisor.service
[Unit]
Description=cAdvisor for monitoring
Requires=docker.service
After=docker.service
[Service]
User=core
TimeoutStartSec=15
ExecStartPre=-/usr/bin/docker pull mastodonc/cadvisor:latest
ExecStart=/usr/bin/docker run --name=cadvisor --label application=cadvisor --net=host --volume=/:/rootfs:ro --volume=/var/run:/var/run:rw --volume=/sys:/sys:ro --volume=/var/lib/docker/:/var/lib/docker:ro --publish=9922:9922  --detach=true mastodonc/cadvisor:latest -port 9922 -storage_driver=influxdb -storage_driver_host=${influxdb-ip}:8086 -storage_driver_db=cadvisor_${cluster-name}
