# resources/systemd/confd.service
[Unit]
Description=running confd
After=multi-user.target
Wants=multi-user.target

[Service]
Restart=always
StandardOutput=journal+console
StandardError=journal+console
ExecStart=/opt/mesosphere/bin/confd -backend=zookeeper -log-level=debug -node=$${internal-lb-dns}:2181 -interval=60
ExecStop=-/usr/bin/pkill confd