# resources/systemd/docker-cleanup.service
[Unit]
Description=Regularly cleanup unused docker images

[Service]
Type=oneshot
ExecStart=/usr/bin/docker rmi $(/usr/bin/docker images -a -q)