# resources/systemd/elastalert-rules.timer
[Unit]
Description=Run Elastalert rule repo github pull every 5minutes

[Timer]
OnCalendar=*:0/2

[Install]
WantedBy=timers.target

