# resources/systemd/var-lib.mount
[Unit]
Description=Mount /var/lib
Before=dbus.service
[Mount]
What=/dev/xvdb
Where=/var/lib
Type=ext4
